# systemd/sol-snatcher.service
[Unit]
Description=Sol $natcher - Autonomous Solana Sniping Bot
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/sol-snatcher
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0
Restart=on-failure
RestartSec=5

# Security settings
User=root
Group=root

# Environment
Environment=COMPOSE_PROJECT_NAME=sol-snatcher

[Install]
WantedBy=multi-user.target